# ops/systemd/deborgen-coordinator.service
[Unit]
Description=deborgen coordinator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace these values for the target host before installing.
User=CHANGE_ME
WorkingDirectory=/absolute/path/to/deborgen
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/deborgen/coordinator.env
# Replace this with the absolute path to uv if it is not on the system service PATH.
ExecStart=/absolute/path/to/uv run deborgen-coordinator
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
